# systemd/hyprland-session-restore.timer
[Unit]
Description=Periodic Hyprland session snapshot

[Timer]
OnActiveSec=60
OnUnitActiveSec=60

[Install]
WantedBy=graphical-session.target
